# config/systemd/user/emacs.service.d/env_vars.conf
[Service]
EnvironmentFile=%h/.profile.d/env.d/user.env
